# as-ap-posang/systemd/ioc-ts-ap-posang.service
[Unit]
Description=Soft IOC for TS Position and Angle Correction
After=multi-user.target

[Service]
Group=sirius
User=sirius
Type=simple
WorkingDirectory=/home/sirius/iocs
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/bash --login -c 'env > /tmp/bashrc-sirius-ioc-ts-ap-posang'
EnvironmentFile=-/tmp/bashrc-sirius-ioc-ts-ap-posang
ExecStart=/usr/local/bin/sirius-ioc-ts-ap-posang.py
