# as-ap-opticscorr/systemd/ioc-bo-ap-tunecorr.service
[Unit]
Description=Soft IOC for BO Tune Correction
After=multi-user.target

[Service]
Group=sirius
User=sirius
Type=simple
WorkingDirectory=/home/sirius/iocs/bo-ap-tunecorr
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/bash --login -c 'env > /tmp/bashrc-sirius-ioc-bo-ap-tunecorr'
EnvironmentFile=-/tmp/bashrc-sirius-ioc-bo-ap-tunecorr
ExecStart=/usr/local/bin/sirius-ioc-bo-ap-tunecorr.py
