# as-ap-opticscorr/systemd/sirius-ioc-bo-ap-chromcorr.service
[Unit]
Description=Soft IOC for BO Chromaticity Correction
After=multi-user.target

[Service]
Group=sirius
User=sirius
Type=simple
WorkingDirectory=/home/sirius/iocs-log/bo-ap-chromcorr
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/bash --login -c 'env > /home/sirius/iocs-log/env-sirius-ioc-bo-ap-chromcorr'
EnvironmentFile=-/home/sirius/iocs-log/env-sirius-ioc-bo-ap-chromcorr
ExecStart=/usr/local/bin/sirius-ioc-bo-ap-chromcorr.py
